# track-studio.service
[Unit]
Description=Track Studio Web App
After=network.target

[Service]
Type=simple
User=benmesch
WorkingDirectory=/home/benmesch/claude_try/track-studio
ExecStart=/usr/bin/python3 /home/benmesch/claude_try/track-studio/app.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
